# step-ca/com.dotfiles.step-ca.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.dotfiles.step-ca</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__DOTFILES_DIR__/step-ca/start.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>STEPPATH</key>
        <string>__STEP_HOME__</string>
        <key>PATH</key>
        <string>__PATH__</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/step-ca.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/step-ca.err.log</string>
</dict>
</plist>
